# systemfiles/crontab
*/5 * * * *  root bash ${NTP_HWCLK_SCRIPT}
